# systemd/rigctld.service
[Unit]
Description=Hamlib rigctld daemon
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# Use a dedicated non-login account when possible; update User= for your system.
User=pi
Group=dialout
EnvironmentFile=-/etc/default/rigctld
Environment=RIGCTLD_MODEL=1042
Environment=RIGCTLD_DEVICE=/dev/ttyUSB0
Environment=RIGCTLD_BAUD=38400
Environment=RIGCTLD_BIND=127.0.0.1
Environment=RIGCTLD_PORT=4532
# -m rig model, -r serial device, -s baud, -T bind address, -t TCP port
ExecStart=/usr/bin/rigctld -m ${RIGCTLD_MODEL} -r ${RIGCTLD_DEVICE} -s ${RIGCTLD_BAUD} -T ${RIGCTLD_BIND} -t ${RIGCTLD_PORT}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
